# contrib/campcli-daemon.service
[Unit]
Description=campcli BC Parks availability daemon
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
EnvironmentFile=%h/.config/campcli/telegram.env
ExecStart=%h/.local/bin/campcli daemon
Restart=on-failure
RestartSec=15
# Keep app logs out of the global journal; write to a dedicated file.
# StateDirectory creates %S/campcli (~/.local/state/campcli) before log setup.
StateDirectory=campcli
StandardOutput=append:%S/campcli/daemon.log
StandardError=append:%S/campcli/daemon.log

[Install]
WantedBy=default.target
