# contrib/campcli-daemon.service
[Unit]
Description=campcli BC Parks availability daemon
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
EnvironmentFile=%h/.config/campcli/telegram.env
ExecStartPre=/usr/bin/mkdir -p %h/.local/state/campcli
ExecStart=%h/.local/bin/campcli daemon
Restart=on-failure
RestartSec=15
# Keep app logs out of the global journal; write to a dedicated file.
StandardOutput=append:%h/.local/state/campcli/daemon.log
StandardError=append:%h/.local/state/campcli/daemon.log

[Install]
WantedBy=default.target
